# Final_Project/03_pipeline_flu_nt/crontab
0 0,12 * * * cd /flu_nt_cronjob && /usr/local/bin/python get_flu_nt.py >> cron.log 2>&1
