# Final_Project/02_pipeline_flu_id/crontab
0 0,6,12,18 * * * cd /flu_id_cronjob && /usr/local/bin/python get_flu_id.py >> cron.log 2>&1
